# quay-harbor-mirror/quay-harbor-mirror.cron
# /etc/cron.d/quay-harbor-mirror
# Mirror release tags of quay.io/ascend/vllm-ascend and vllm/vllm-openai into the
# local Harbor every 30 min. flock -n skips a cycle if a previous run is still busy
# (large image pulls can take longer than 30 min).
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
7,37 * * * * root flock -n /var/lock/quay-harbor-mirror.lock /usr/local/bin/quay-harbor-mirror.sh >> /var/log/quay-harbor-mirror.log 2>&1
